# docker/bastion/crontab
# run checkcerts.sh every minute
*/1  *    * * *   root /usr/local/bin/checkcerts.sh 1>/proc/1/fd/1 2>/proc/1/fd/2
